# motd-trigger.service
[Unit]
Description=Update message of the day
BindsTo=update-triggers.target

[Service]
Type=oneshot
ExecStart=/usr/bin/motd-update.sh

